# crossGameConfBase/etc/rc.local
#!/bin/bash

if chmod 777 /dev/gpiomem; then
    echo "Added privileges to /dev/gpiomem"
else
    echo "Error adding privileges to /dev/gpiomem"
    exit 1
fi

echo "Waiting for ssd..."
while ! ls "/rpi/resources/" ;do
    sleep 0.5
done

systemctl start smbd.service

exit 0